# base/conf/supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log
logfile_backups=10
logfile_maxbytes=50MB
loglevel=info
pidfile=/var/run/supervisord.pid

[program:sshd]
command=/usr/sbin/sshd -D
pidfile=/var/run/sshd.pid
stdout_logfile=/var/log/supervisor/%(program_name)s.stdout
stderr_logfile=/var/log/supervisor/%(program_name)s.stderr

[include]
files=/etc/supervisord.d/*.conf
